# templates/claude-project-manager.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>claude-project-manager</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>{{PROJECT_DIR}}/cpm</string>
        <string>run</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/claude-project-manager/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/claude-project-manager/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{PATH}}</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>{{PROJECT_DIR}}</string>
</dict>
</plist>
